# debian/vboxadm-sa.prerm
#!/bin/sh -e

# VBoxAdm SA prerm
# Dominik Schulz <dominik.schulz@gauner.org>

case "$1" in
  upgrade)
    new=$2 # new version
    invoke-rc.d --quiet vboxadm stop
    ;;

  deconfigure)
    ;;

  remove)
    invoke-rc.d --quiet vboxadm stop
    ;;

  failed-upgrade)
    ;;

  *)
    echo "prerm called with unknown argument \`$1'" >&2
    exit 1
    ;;
esac
#DEBHELPER#
exit 0
